# tickets-watcher.service
# systemd unit для tickets-watcher.
#
# Поправь под себя:
#   - User=            — под каким пользователем запускать
#   - WorkingDirectory — куда положил проект
#   - ExecStart        — путь к python из venv проекта
#
# Установка:
#   sudo cp tickets-watcher.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tickets-watcher
# Логи:
#   journalctl -u tickets-watcher -f    (плюс watcher.log в папке проекта)

[Unit]
Description=Yandex Afisha ticket watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=youruser
WorkingDirectory=/home/youruser/tickets-watcher
ExecStart=/home/youruser/tickets-watcher/venv/bin/python tickets_watcher.py
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
